# systemd/pantalla-kiosk-watchdog@.service.d/10-rollback.conf
[Unit]
ConditionPathExists=!/var/lib/pantalla-reloj/state/disable-kiosk-watchdog
