# systemd/pantalla-portal@.service
[Unit]
Description=Pantalla_reloj XDG desktop portal for user %i (optional)
After=pantalla-openbox@%i.service

[Service]
User=%i
PermissionsStartOnly=yes
Environment=XDG_RUNTIME_DIR=/run/user/%U
Environment=DISPLAY=:0
Environment=XAUTHORITY=/home/%i/.Xauthority
EnvironmentFile=-/var/lib/pantalla-reloj/state/session.env
ExecStartPre=/usr/bin/test -x /usr/libexec/xdg-desktop-portal-gtk
ExecStartPre=/usr/bin/test -x /usr/libexec/xdg-desktop-portal
ExecStartPre=/bin/sh -c 'test -d /var/log/pantalla && test -w /var/log/pantalla'
ExecStartPre=/bin/sh -c 'pkill -u %U -f "xdg-desktop-portal" >/dev/null 2>&1 || true'
ExecStartPre=/bin/sh -c 'pkill -u %U -f "xdg-desktop-portal-gtk" >/dev/null 2>&1 || true'
ExecStart=/opt/pantalla/bin/pantalla-portal-launch.sh
Restart=on-failure
RestartSec=2

[Install]
WantedBy=multi-user.target
