# systemd/pantalla-kiosk@.service.d/10-sanitize-rollback.conf
[Service]
Environment=PANTALLA_SANITIZER_DISABLE_FILE=/var/lib/pantalla-reloj/state/disable-kiosk-sanitize
